# deploy/arb-bot.service
[Unit]
Description=Arb Bot Config Service
After=network.target

[Service]
WorkingDirectory=/opt/arb-bot
User=deploy
Group=deploy
EnvironmentFile=/opt/arb-bot/.env
ExecStart=/opt/arb-bot/.venv/bin/uvicorn services.config_service.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
